# deploy/systemd/light-claw.service
[Unit]
Description=light-claw multi-agent bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wysi
WorkingDirectory=/data1/wysi/light-claw
EnvironmentFile=/data1/wysi/light-claw/.env
Environment="PATH=/data2/wysi/miniconda3/bin:/home/wysi/.bun/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/data2/wysi/miniconda3/bin/uv run light-claw
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
